# roles/ufw/templates/ufw-docker.service.j2
# SPDX-License-Identifier: MIT-0
# UFW Docker integration service
# Configures DOCKER-USER chain after Docker starts
# Managed by Ansible - do not edit manually

[Unit]
Description=UFW Docker integration (DOCKER-USER chain)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ufw-docker-rules.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
