# roles/packages/templates/apt_update.timer.j2
[Unit]
Description=Run apt update & upgrade daily

[Timer]
OnCalendar={{ apt_update_timer_on_calendar }}
Persistent={{ apt_update_timer_persistent | ternary('true','false') }}

[Install]
WantedBy=timers.target
